# gobgp/gobgpd-network.timer
[Unit]
Description=GoBGP Routing Network Inject Timer

[Timer]
OnActiveSec=30s

[Install]
WantedBy=timers.target